# setx11locale.service
[Unit]
Description=Set the X11 keyboard layout to german

[Service]
ExecStart=/usr/bin/setx11locale

[Install]
WantedBy=multi-user.target 
